# systemd/litclock-all.service
[Unit]
Description=Literary Clock with Web Interface
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/litclock
ExecStart=/home/pi/litclock/venv/bin/python -m litclock.cli
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=litclock-all

[Install]
WantedBy=multi-user.target 